# installer/supervisor.d/dev/mobile-wallet-service.conf
[program:reservation-service]
command                 = /home/vagrant/.nvm/nvm-exec node index.js
process_name            = mobile-wallet-service
user                    = vagrant
numprocs                = 1
autostart               = true
autorestart             = true
directory               = /var/www/mobile-wallet-service
environment             = NODE_ENV=dev
